# .config/systemd/user/vdirsyncer.service
[Unit]
Description=Synchronize calendars and contacts via vdirsyncer

[Service]
ExecStart=/usr/bin/vdirsyncer sync

